# deploy/systemd/titanflow-mod@.service
[Unit]
Description=TitanFlow Module: %i
After=titanflow-core.service
BindsTo=titanflow-core.service
PartOf=titanflow-core.service

[Service]
Type=simple
User=titanflow
Group=titanflow

ExecStart=/opt/titanflow/.venv/bin/python -m titanflow.modules.%i.module
WorkingDirectory=/opt/titanflow

ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateUsers=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service

RestrictAddressFamilies=AF_UNIX

RuntimeDirectory=titanflow-mod-%i
RuntimeDirectoryMode=0700
ReadWritePaths=/run/titanflow-mod-%i
ReadOnlyPaths=/etc/titanflow/manifests/%i.yaml

MemoryMax=512M
CPUWeight=50
IOWeight=50

Environment=TITANFLOW_MODULE_ID=%i
Environment=TITANFLOW_CORE_SOCKET=/run/titanflow/core.sock
Environment=TITANFLOW_CONFIG_DIR=/opt/titanflow/config
EnvironmentFile=-/etc/titanflow/secrets/%i.env

Restart=on-failure
RestartSec=10

[Install]
WantedBy=titanflow-core.service
