# deploy/systemd/titanflow-v03-telemetry-http.service
[Unit]
Description=TitanFlow v0.3 Telemetry HTTP Bridge
After=network-online.target titanflow-v03.service
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
RuntimeDirectory=titanflow-telemetry
WorkingDirectory=/opt/titanflow
ExecStart=/opt/titanflow/.venv/bin/python -m titanflow.v03.telemetry_http_runner

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
