# deploy/systemd/titanflow-v03-telemetry.socket
[Unit]
Description=TitanFlow v0.3 Telemetry Socket

[Socket]
ListenStream=/run/titanflow-core/telemetry.sock
SocketMode=0600
Backlog=64

[Install]
WantedBy=sockets.target
