# deploy/systemd/titanflow-v03.service
[Unit]
Description=TitanFlow v0.3 Microkernel
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
DynamicUser=yes
StateDirectory=titanflow/core
RuntimeDirectory=titanflow-core
RuntimeDirectoryMode=0750

WorkingDirectory=/opt/titanflow
ExecStart=/opt/titanflow/.venv/bin/python -m titanflow.v03.run

Environment=TITANFLOW_DB_PATH=/var/lib/titanflow/core/titanflow.db
Environment=TITANFLOW_TELEMETRY_SOCKET=/run/titanflow-core/telemetry.sock
Environment=TITANFLOW_CORE_SOCKET=/run/titanflow-core/core.sock

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/titanflow/core /run/titanflow-core
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

Restart=always
RestartSec=5
WatchdogSec=20

[Install]
WantedBy=multi-user.target
